# systemd/.config/systemd/user/onedrive.service
[Unit]
Description=Onedrive

[Service]
Type=simple
ExecStart=/usr/bin/onedrive --monitor

[Install]
WantedBy=default.target
